# services/dropbox-daemon.service
[Unit]
Description=headless Dropbox
RefuseManualStart=no

[Service]
Type=oneshot
ExecStart=/srv/dropbox/daemon.sh

